# serving/contrib/apprentice-vllm@.service
[Unit]
Description=Apprentice vLLM serving for pattern %i
After=network-online.target registry-service.service
Wants=network-online.target
Requires=registry-service.service

[Service]
Type=simple
# Override User= with `systemctl edit apprentice-vllm@<pattern>` to run as
# your account, or create a dedicated `apprentice` user that owns
# ~/.apprentice/registry.  Do NOT leave User= unset — vLLM as root has full
# access to the host's GPU device files with no isolation.
User=apprentice
Group=apprentice
ExecStartPre=/usr/bin/apprentice-serve --check-only --pattern-id %i
ExecStart=/usr/bin/apprentice-serve --pattern-id %i
Restart=on-failure
RestartSec=10
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.apprentice/registry
PrivateTmp=yes

[Install]
WantedBy=default.target
